# .config/systemd/user/checkupdates.service
[Unit]
Description=Check for updates every hour

[Service]
ExecStart=%h/.dwm/bar-scripts/updates.sh

