# .config/systemd/user/cpu_usage.service
[Unit]
Description=CPU Usage Monitor

[Service]
ExecStart=/home/cie/.bin/cpu_usage.sh
