# .config/systemd/user/cpu_temp.service
[Unit]
Description=CPU Temp Monitor

[Service]
ExecStart=/home/cie/.bin/cpu_temp.sh
