# .config/systemd/user/cpu_temp.timer
[Unit]
Description=Run CPU Temp Monitor every 5 seconds

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
AccuracySec=1s

[Install]
WantedBy=timers.target
