# .config/systemd/user/disk_space.service
[Unit]
Description=Disk Space Monitor

[Service]
ExecStart=/home/cie/.bin/disk_space.sh
