# .config/systemd/user/dwm-bar.service
[Unit]
Description=dwm-bar

[Service]
Environment=DISPLAY=:0
ExecStart=/home/cie/.config/systemd/user/scripts/dwm-bar.sh
