# .config/systemd/user/dwm-bar.service
[Unit]
Description=dwm-bar

[Service]
Environment=DISPLAY=:0
ExecStart=/home/cie/.bin/dwm-bar.sh
