# .config/systemd/user/easy_active.timer
[Unit]
Description=Timer unit for the easy_active service

[Timer]
OnBootSec=5s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
