# .config/systemd/user/mem_usage.timer
[Unit]
Description=Run Mem Usage Monitor every 2 seconds

[Timer]
OnBootSec=2s
OnUnitActiveSec=2s
AccuracySec=1s

[Install]
WantedBy=timers.target
