# .config/systemd/user/mem_usage.service
[Unit]
Description=Mem Usage Monitor

[Service]
ExecStart=/home/cie/.bin/mem_usage.sh
