# .config/systemd/user/music.service
[Unit]
Description=The script to switch between bass and equalizer.

[Service]
ExecStart=/home/cie/.bin/music.sh
