# .config/systemd/user/pacman_info.timer
[Unit]
Description=Timer to run the pacman_info script

[Timer]
OnBootSec=5s
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
