# .config/systemd/user/santos_matchday.service
[Unit]
Description=Santos matchday script

[Service]
ExecStart=/home/cie/.bin/santos_matchday.py
