# .config/systemd/user/santos_calender.service
[Unit]
Description=A script to display Santos fc matchs

[Service]
Type=simple
ExecStart=/home/digo/.bin/santos_calendar.mjs
