# .config/systemd/user/santos_schedule.service
[Unit]
Description=Santos schedule script

[Service]
ExecStart=/home/cie/.bin/santos_schedule.py
