# .config/systemd/user/server_check.service
[Unit]
Description=check the server info to display

[Service]
ExecStart=%h/.bin/server_check.sh
WorkingDirectory=/home/cie/.bin/

