# .config/systemd/user/server.timer
[Unit]
Description=Checks for TF2 servers

[Timer]
OnCalendar=*-*-* 05:00/30:00
Persistent=true

[Install]
WantedBy=timers.target
