# .config/systemd/user/server.service
[Unit]
Description=tf2 server

[Service]
ExecStart=/home/cie/.bin/server_info.sh
