# .config/systemd/user/tf2_server_display.timer
[Unit]
Description=Timer to Update TF2 Player Display Every 30 Seconds

[Timer]
OnBootSec=5s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
