# .config/systemd/user/tf2_server_count.timer
[Unit]
Description=Timer to Fetch TF2 Player Counts

[Timer]
OnBootSec=5s
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
